# deploy/systemd/pro-botti-meta-train.service
[Unit]
Description=pro_botti META ensemble training (all symbols)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/pro_botti
EnvironmentFile=/root/pro_botti/secrets.env
EnvironmentFile=/root/pro_botti/botti.env
# ASETTAA symbols.txt:n eksplisiittisesti. VAIHDA polkuun jossa sinun tiedostosi on.
Environment=META_SYMBOLS_FILE=/root/pro_botti/symbols.txt
# Environment=META_TFS=15m,1h,4h
# Environment=META_PARALLEL=6
# Environment=META_MIN_CANDLES=300
# Environment=META_MAX_SYMBOLS=0
ExecStart=/root/pro_botti/venv/bin/python -m tools.meta_train_all
Nice=10
IOSchedulingClass=2
IOSchedulingPriority=7
