# deploy/systemd/pro-botti-meta-train.timer
[Unit]
Description=Run pro_botti META training periodically

[Timer]
OnBootSec=5m
OnUnitActiveSec=30m
RandomizedDelaySec=2m
Persistent=true
Unit=pro-botti-meta-train.service

[Install]
WantedBy=timers.target
