# deploy/systemd/pro_botti-fix-ellipsis.timer
[Unit]
Description=CapitalBot - Weekly Ellipsis Auto-Fix (Sundays 05:00)

[Timer]
OnCalendar=Sun *-*-* 05:00:00
Persistent=true

[Install]
WantedBy=timers.target
