# deploy/systemd/pro_botti-equity.timer
[Unit]
Description=CapitalBot - Daily Equity Curve (22:00 Helsinki)

[Timer]
OnCalendar=*-*-* 19:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
