# deploy/systemd/pro_botti-pnl-monthly.service
[Unit]
Description=CapitalBot - Monthly PnL Summary (1st day 21:00)
After=network.target

[Service]
Type=oneshot
ExecStart=/root/pro_botti/venv/bin/python /root/pro_botti/tools/pnl_report.py 30 monthly
EnvironmentFile=/root/pro_botti/secrets.env
StandardOutput=append:/root/pro_botti/logs/pnl_report.log
StandardError=append:/root/pro_botti/logs/pnl_report.log
WorkingDirectory=/root/pro_botti
