# deploy/systemd/pro_botti-pnl-weekly.timer
[Unit]
Description=CapitalBot - Weekly Telegram report (Sunday 21:00 Helsinki)

[Timer]
OnCalendar=Sun *-*-* 18:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
