# deploy/systemd/pro_botti-ppo-trigger.timer
[Unit]
Description=CapitalBot PPO retrain trigger timer (check drift/PnL hourly)

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
