# deploy/pro_botti.service
[Unit]
Description=Pro Botti - Live Trading Service
Documentation=https://github.com/mikkotorma82-lgtm/pro_botti
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/pro_botti
Environment="ROOT=/root/pro_botti"
Environment="PATH=/root/pro_botti/venv/bin:/usr/local/bin:/usr/bin:/bin"

# Load environment variables from botti.env
EnvironmentFile=-/root/pro_botti/botti.env

# Start live trading with active symbols
ExecStart=/root/pro_botti/venv/bin/python3 -m cli live --config /root/pro_botti/config.yaml

# Logging
StandardOutput=append:/root/pro_botti/logs/live.log
StandardError=append:/root/pro_botti/logs/live.err

# Restart policy
Restart=always
RestartSec=10

# Resource limits
MemoryMax=2G
CPUQuota=100%

# Signal handling
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
